# conf/supervisor/chop_chop.conf
; ==================================
;  celery worker supervisor example
; ==================================

; the name of your supervisord program
[program:chop_chop]

; Set full path to celery program if using virtualenv
command=/home/interdev/app/venv/bin/celery worker -A backend

; Set the enviroment variables
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,SERVER_NAME="PROD_SERVER"

; The directory to your Django project
directory=/home/interdev/app/backend

; If supervisord is run as the root user, switch users to this UNIX user account
; before doing any processing.
user=interdev

; Supervisor will start as many instances of this program as named by numprocs
numprocs=1

; Put process stdout output in this file
stdout_logfile=/home/interdev/app/logs/chop_chop.log

; Put process stderr output in this file
stderr_logfile=/home/interdev/app/logs/chop_chop_err.log

; If true, this program will start automatically when supervisord is started
autostart=true


